# 20131108.crontab
0 11 * * mon,wed,fri /home/tbg/work/bin/remind.sh
0 23 * * mon,wed,fri /home/tbg/work/bin/runtbg.sh
